# deploy/tachtach-kiosk.service
[Unit]
Description=TachTach-Screens Chromium Kiosk
After=tachtach-server.service graphical.target
Wants=tachtach-server.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# NOTE: __PI_USER__ / __PI_HOME__ are placeholder tokens. This file is a
# template — pi-setup.sh substitutes the real username and home directory
# before copying it into /etc/systemd/system/. Do not install this file
# verbatim; run pi-setup.sh instead.
User=__PI_USER__
Group=__PI_USER__
Environment=DISPLAY=:0
Environment=XAUTHORITY=__PI_HOME__/.Xauthority

# Wait for server to be ready, then launch kiosk
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/chromium-browser \
    --kiosk \
    --disable-restore-session-state \
    --noerrdialogs \
    --disable-infobars \
    --disable-translate \
    --no-first-run \
    --start-fullscreen \
    --check-for-update-interval=31536000 \
    --overscroll-history-navigation=0 \
    --disable-pinch \
    --disable-features=TranslateUI \
    --disable-dev-shm-usage \
    --autoplay-policy=no-user-gesture-required \
    http://127.0.0.1:3000/screen

Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tachtach-kiosk

[Install]
WantedBy=graphical.target
